# clawmander-backend.service
[Unit]
Description=Clawmander Dashboard Backend
After=network.target

[Service]
Type=simple
Environment="PATH=%h/.npm-global/bin:%h/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"
WorkingDirectory=%h/sandbox/personal/clawmander/backend
ExecStart=/home/linuxbrew/.linuxbrew/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
